# flight-bot.service
[Unit]
Description=Flight Monitor Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/projects/flight-monitor
Environment=TELEGRAM_BOT_TOKEN=8537390248:AAGsxHWWdV3siXrk9-xgAwGu5g6VNzmFp0k
Environment=TELEGRAM_CHAT_ID=137259321
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/projects/flight-monitor/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
